# packages/freesmartphone/python-ousaged/ousaged
#! /bin/sh
#
# odeviced  This shell script starts and stops the open device daemon.
#
# chkconfig: 345 90 20
# description: py-ousaged is the open usage daemon
# processname: python

PATH=/bin:/usr/bin:/sbin:/usr/sbin
NAME=odeviced

[ -f /etc/default/rcS ] && . /etc/default/rcS

case "$1" in
    start)
        echo -n "Starting open usage daemon: "
        start-stop-daemon --start --pidfile /var/run/${NAME}.pid --make-pidfile --background -x /usr/bin/ousaged
        if [ $? = 0 ]; then
            echo "(ok)"
        else
            echo "(failed)"
        fi
        ;;
    stop)
        echo -n "Stopping open usage daemon: "
        start-stop-daemon --stop --pidfile /var/run/${NAME}.pid --oknodo
        rm -f /var/run/${NAME}.pid
        echo "(done)"
        ;;
    restart|force-reload)
        $0 stop
        $0 start
        ;;
    *)
        echo "Usage: /etc/init.d/ousaged {start|stop|restart|force-reload}"
        exit 1
        ;;
esac

exit 0
