# recipes/images/shr-image/etc/profile.d/elementary.sh
export ELM_ENGINE=x11
export ELM_THEME=default
export ELM_SCALE=2
export ELM_FINGER_SIZE=70
